# supervisord.conf
[supervisord]
nodaemon = true

[program:api]
command = pipenv run api
autorestart = true
stdout_events_enabled=true
stderr_events_enabled=true
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
